# src/deb/init/ned.upstart
description "Named Entity Database Service"

start on started and runlevel [2345]
stop on shutdown

respawn
respawn limit 5 15

script
exec /bin/bash << EOT
  export CATALINA_BASE=/opt/plos/ned
  /usr/share/tomcat7/bin/catalina.sh run
EOT
end script
